# dist/linux/tagent.service
[Unit]
Description=TrustAgent Service

[Service]
Type=simple
User=tagent
Group=tagent
WorkingDirectory=/opt/trustagent/bin
ExecStart=/opt/trustagent/bin/tagent start
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStartSec=15
Restart=on-failure
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target