# dist/linux/tagent.service
[Unit]
Description=TrustAgent Service

[Service]
Type=forking
WorkingDirectory=/opt/trustagent/bin
ExecStart=/opt/trustagent/bin/tagent start
ExecStop=/opt/trustagent/bin/tagent stop
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStartSec=45
Restart=on-failure
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target